# tests/modules/services/wob/wob-service-expected.socket
[Install]
WantedBy=sockets.target

[Socket]
FlushPending=yes
ListenFIFO=%t/wob.sock
RemoveOnStop=yes
SocketMode=0600
