# tests/modules/services/tldr-update/tldr-update.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=monthly
Persistent=true

[Unit]
Description=Update tldr CLI cache
