# tests/modules/services/nix-gc/expected.service
[Service]
ExecStart=@nix@/bin/nix-collect-garbage --delete-older-than 30d

[Unit]
Description=Nix Garbage Collector
