# tests/modules/services/nix-gc/expected.service
[Service]
ExecStart=/nix/store/00000000000000000000000000000000-nix-gc
Type=oneshot

[Unit]
Description=Nix Garbage Collector
