# tests/modules/services/lieer/lieer-service-expected.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*:0/5
RandomizedDelaySec=30

[Unit]
Description=lieer Gmail synchronization for hm@example.com
