# tests/modules/services/podman/linux/build-expected.service
# Automatically generated by /nix/store/00000000000000000000000000000000-podman/lib/systemd/user-generators/podman-user-generator
# 
# Automatically generated by home-manager for podman build configuration
# DO NOT EDIT THIS FILE DIRECTLY
#
# my-bld.build
[X-Build]
Environment=
File=/nix/store/00000000000000000000000000000000-Containerfile
ImageTag=homemanager/my-bld
Label=nix.home-manager.managed=true
TLSVerify=true

[Install]
WantedBy=default.target
WantedBy=multi-user.target

[Service]
RemainAfterExit=yes
TimeoutStartSec=300
ExecStart=/nix/store/00000000000000000000000000000000-podman/bin/podman build --tls-verify --tag homemanager/my-bld --label nix.home-manager.managed=true --file /nix/store/00000000000000000000000000000000-Containerfile
SyslogIdentifier=%N
Type=oneshot

[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Description=Service for build my-bld
RequiresMountsFor=%t/containers
SourcePath=/nix/store/00000000000000000000000000000000-home-build-podman-my-bld/quadlets/podman-my-bld.build
