# tests/modules/services/podman/linux/network-expected.service
# Automatically generated by /nix/store/00000000000000000000000000000000-podman/lib/systemd/user-generators/podman-user-generator
# 
# Automatically generated by home-manager for podman network configuration
# DO NOT EDIT THIS FILE DIRECTLY
#
# my-net.network
[Install]
WantedBy=default.target
WantedBy=multi-user.target

[X-Network]
Gateway=192.168.1.1
Label=nix.home-manager.managed=true
NetworkName=my-net
Options=isolate=true
PodmanArgs=--ipam-driver dhcp
PodmanArgs=--dns=192.168.55.1
PodmanArgs=--log-level=debug
Subnet=192.168.1.0/24

[Service]
Environment=PATH=/run/wrappers/bin:/usr/bin:/bin:/usr/sbin:/sbin:@shadow@/bin:/nix/store/00000000000000000000000000000000-coreutils/bin
ExecStartPre=/nix/store/00000000000000000000000000000000-await-podman-unshare
RemainAfterExit=yes
TimeoutStartSec=15
ExecStart=/nix/store/00000000000000000000000000000000-podman/bin/podman network create --ignore --subnet 192.168.1.0/24 --gateway 192.168.1.1 --opt isolate=true --label nix.home-manager.managed=true --ipam-driver dhcp --dns=192.168.55.1 --log-level=debug my-net
SyslogIdentifier=%N
Type=oneshot

[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
After=network.target
Description=Service for network my-net
SourcePath=/nix/store/00000000000000000000000000000000-home-network-podman-my-net/quadlets/podman-my-net.network
RequiresMountsFor=%t/containers
