# tests/modules/misc/xsession/basic-setxkbmap-expected.service
[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=@setxkbmap@/bin/setxkbmap -layout 'us' -variant ''
RemainAfterExit=true
Type=oneshot

[Unit]
After=graphical-session-pre.target
Description=Set up keyboard in X
PartOf=graphical-session.target
