# tests/modules/misc/xsession/keyboard-without-layout-expected.service
[Install]
WantedBy=graphical-session.target

[Service]
ExecStart=@setxkbmap@/bin/setxkbmap -option 'ctrl:nocaps' -option 'altwin:no_win'
RemainAfterExit=true
Type=oneshot

[Unit]
After=graphical-session-pre.target
Description=Set up keyboard in X
PartOf=graphical-session.target
