# tests/modules/services/borgmatic/linux/basic-configuration.timer
[Install]
WantedBy=timers.target

[Timer]
OnCalendar=weekly
Persistent=true
RandomizedDelaySec=10m

[Unit]
Description=Run borgmatic backup
