# tests/modules/services/podman-linux/volume-expected.service
# Automatically generated by /nix/store/00000000000000000000000000000000-podman/lib/systemd/user-generators/podman-user-generator
# 
# Automatically generated by home-manager for podman volume configuration
# DO NOT EDIT THIS FILE DIRECTLY
#
# my-vol.volume
[Install]
WantedBy=default.target
WantedBy=multi-user.target

[Service]
Environment=PATH=/run/wrappers/bin:/usr/bin:/bin:/usr/sbin:/sbin:@shadow@/bin:/nix/store/00000000000000000000000000000000-coreutils/bin
ExecStartPre=/nix/store/00000000000000000000000000000000-await-podman-unshare
RemainAfterExit=yes
TimeoutStartSec=15
ExecStart=/nix/store/00000000000000000000000000000000-podman/bin/podman volume create --ignore --opt copy --opt device=tmpfs --opt type=tmpfs --opt o=uid=1000,gid=1000 --label nix.home-manager.managed=true --label nix.home-manager.preserve=true --module=/etc/nvd.conf my-vol
SyslogIdentifier=%N
Type=oneshot

[Unit]
Wants=podman-user-wait-network-online.service
After=podman-user-wait-network-online.service
Description=Service for volume my-vol
SourcePath=/nix/store/00000000000000000000000000000000-home-volume-podman-my-vol/quadlets/podman-my-vol.volume
RequiresMountsFor=%t/containers

[X-Volume]
Copy=true
Device=tmpfs
Group=1000
Label=nix.home-manager.managed=true
Label=nix.home-manager.preserve=true
PodmanArgs=--module=/etc/nvd.conf
Type=tmpfs
User=1000
VolumeName=my-vol
