# tests/modules/programs/beets/mpdstats-external-expected.service
[Install]
WantedBy=default.target

[Service]
ExecStart=@beets@/bin/beet mpdstats

[Unit]
Description=Beets MPDStats daemon
