# tests/modules/programs/opencode/web-service-environment-file.service
[Install]
WantedBy=default.target

[Service]
Environment=PATH=/home/hm-user/.nix-profile/bin
EnvironmentFile=/run/secrets/opencode
ExecStart=@opencode@/bin/opencode serve 
Restart=always
RestartSec=5

[Unit]
After=network.target
Description=OpenCode Web Service
