# tests/modules/services-modular/demo.service
[Install]
WantedBy=default.target

[Service]
ExecStart="@mpd@/bin/mpd" "--no-daemon"
Restart=always
RestartSec=5
Type=simple
